# .config/systemd/user/battery-notifier.timer
[Unit]
Description=Battery Notification Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=1min

[Install]
WantedBy=default.target

