# .config/systemd/user/check_power.service
[Unit]
Description=Check Power Status

[Service]
Type=simple
ExecStart=/home/neon/.local/bin/check_power.sh

[Install]
WantedBy=default.target
